# sleep-on-lan.conf
start on (net-device-up
          and local-filesystems
	  and runlevel [2345])
stop on runlevel [016]

respawn

chdir /opt/sleep-on-lan
exec /opt/sleep-on-lan/ruby_env.sh /opt/sleep-on-lan/sol.rb
